# qmd/com.tobilu.qmd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tobilu.qmd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__QMD_DIST__</string>
        <string>mcp</string>
        <string>--http</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/qmd-mcp.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/qmd-mcp.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
