# qmd/qmd-mcp.service
[Unit]
Description=qmd MCP HTTP server
After=network.target

[Service]
Type=simple
ExecStart=__BUN_PATH__ __QMD_SRC__ mcp --http
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
